# deploy/systemd/medgreenrev.docker-compose.service.template
[Unit]
Description=MEDGREENREV app (Docker Compose)
After=network.target docker.service
Requires=docker.service

[Service]
User=${DOCKER_USER}
WorkingDirectory=${DOCKER_COMPOSE_DIRECTORY}
ExecStart=${DOCKER_COMPOSE_DIRECTORY}/deploy/up.sh
Restart=always
RestartSec=10
TimeoutStartSec=0

StandardOutput=journal
StandardError=journal


[Install]
WantedBy=multi-user.target
